# whisper-ptt.service
[Unit]
Description=Whisper Push-to-Talk (Right Ctrl)
After=graphical-session.target pipewire.service
Wants=pipewire.service

[Service]
Type=simple
Environment=PYTHONUNBUFFERED=1
# LD_LIBRARY_PATH is set dynamically by install.sh with resolved Python version
ExecStart=%h/.local/bin/whisper-ptt
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
